# demo-service/supervisord.conf
[supervisord]
nodaemon=true

[program:demo-service]
command=/app/demo-service

#不需要单独提供 HTTP 服务端口，可借助 Go Micro 内置的 API 网关提供
#[program:http]
#command=/app/demo-service -mode=http

